# ssh-terminal-bot.service
[Unit]
Description=SSH Terminal Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ssh-terminal
EnvironmentFile=/opt/ssh-terminal/.env
ExecStart=/usr/bin/python3 /opt/ssh-terminal/bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
